# binary-creator/supervisor.conf
[supervisord]
nodaemon=true

[program:map-renderer]
command=nodejs /app
killasgroup=true
stopasgroup=true
redirect_stderr=true

[program:chromium]
command=/usr/bin/chromium --headless --no-sandbox --disable-gpu --disable-translate --disable-dev-shm-usage --disable-extensions --disable-background-networking --safebrowsing-disable-auto-update --disable-sync --metrics-recording-only --disable-default-apps --no-first-run--mute-audio --hide-scrollbars --remote-debugging-port=9222
killasgroup=true
stopasgroup=true
redirect_stderr=true